# templates/upstart/process.conf.erb
start on starting <%= app %>-<%= name %>
stop on stopping <%= app %>-<%= name %>
respawn

exec su -l -s /bin/bash -c 'cd <%= engine.root %>; export PORT=<%= port %>;<% engine.env.each_pair do |var,env| %> export <%= var %>=<%= shell_quote(env) %>; <% end %> exec "$0" "$@"' <%= user %> -- <%= process.command %> >> <%= log %>/<%=name%>-<%=num%>.log 2>&1
